# deploy/feishu-claude.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zara.feishu-claude</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>-u</string>
        <string>/Users/bytedance/Documents/Claude code projects/feishu-claude-code/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/bytedance/Documents/Claude code projects/feishu-claude-code</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/bytedance/.npm-global/bin:/Users/bytedance/.local/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/bytedance</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/bytedance/Documents/Claude code projects/feishu-claude-code/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bytedance/Documents/Claude code projects/feishu-claude-code/stderr.log</string>
</dict>
</plist>
